# systemd/ha-api.service
[Unit]
Description=Home Automation — Query API
Documentation=file:///home/visko/home_automation/home-automation-architecture-plan.md
After=ha-writer.service network.target

[Service]
Type=simple
User=visko
Group=visko
WorkingDirectory=/home/visko/home_automation
Environment=HA_DB=instance/db/hot.db
Environment=HA_PARQUET_DIR=instance/db/parquet
Environment=HA_MAX_DEEP_ROWS=50000
# DuckDB scratch needs a writable home; service home is read-only (ProtectHome).
# instance/ is in ReadWritePaths.
Environment=HOME=/home/visko/home_automation/instance
# Bind address: 0.0.0.0 = reachable on the LAN; 127.0.0.1 = localhost only.
# API is read-only and unauthenticated — only expose on a trusted home LAN.
Environment=HA_BIND_HOST=0.0.0.0
ExecStart=/home/visko/home_automation/venv/bin/python3 \
    -m uvicorn server.api.main:app \
    --host ${HA_BIND_HOST} \
    --port 8123 \
    --workers 1 \
    --no-access-log
Restart=on-failure
RestartSec=10

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/visko/home_automation/instance
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ha-api

[Install]
WantedBy=multi-user.target
