# systemd/ha-cluster-heartbeat.service
[Unit]
Description=Home Automation — Cluster Heartbeat (ha/cluster/<node>/heartbeat for failover sensing)
After=mosquitto.service network.target
Wants=mosquitto.service

[Service]
Type=simple
User=visko
Group=visko
WorkingDirectory=/home/visko/home_automation
Environment=HA_BROKER=localhost
Environment=HA_HEARTBEAT_S=3
EnvironmentFile=-/home/visko/home_automation/instance/mqtt.env
ExecStart=/home/visko/home_automation/venv/bin/python3 -m server.cluster.heartbeat

Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/visko/home_automation/instance
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ha-cluster-heartbeat

[Install]
WantedBy=multi-user.target
