# systemd/ha-tasmota-bridge.service
[Unit]
Description=Home Automation — Tasmota Bridge (WiFi plug/meter → canonical /state)
After=mosquitto.service network.target
Wants=mosquitto.service

[Service]
Type=simple
User=visko
Group=visko
WorkingDirectory=/home/visko/home_automation
Environment=HA_BROKER=localhost
Environment=HA_BROKER_PORT=1883
# Broker auth cutover: dictator MQTT creds (gitignored, optional via '-'; latent on anon broker)
EnvironmentFile=-/home/visko/home_automation/instance/mqtt.env
ExecStart=/home/visko/home_automation/venv/bin/python3 \
    -m server.ingest.tasmota_bridge \
    --registry instance/tasmota-devices.yaml \
    --broker ${HA_BROKER} \
    --broker-port ${HA_BROKER_PORT}
Restart=on-failure
RestartSec=10

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/visko/home_automation/instance
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ha-tasmota-bridge

[Install]
WantedBy=multi-user.target
